# 8.1-nginx/rootfs/etc/s6-overlay/s6-rc.d/supercronic/run
#!/command/execlineb -P
with-contenv
supercronic /etc/supercronic/crontab
